# assets/veloren-server/veloren@.timer
[Unit]
Description=Check for updates periodically

[Timer]
Unit=oneshot-veloren@%i.service
OnCalendar=*:0/15

[Install]
WantedBy=timers.target
